# clash-supervisord.conf
[program:clash]
command=/usr/local/bin/clash -d /clash/data
autostart=true
autorestart=true
user=root
